# debian/postinst
#!/bin/sh

systemctl unmask shellhub.service
systemctl restart shellhub.service
systemctl enable shellhub.service

exit 0
